# docker/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[program:bitcoind]
command=/bin/sh -c "/bin/rm -fr ~bitcoin/.bitcoin/regtest/ && /usr/local/bin/bitcoind -regtest -rpcbind=127.0.0.1 -rpcallowip=127.0.0.1/24 -rpcuser=alice -rpcpassword=password "
directory=/usr/local/bin
environment=HOME="/home/bitcoin",PATH="/usr/bin:/usr/sbin"
user=bitcoin
group=bitcoin
autostart=true
autorestart=true
startsecs=60
startretries=3
stopwaitsecs=600
stopsignal=TERM


[program:fulcrum]
command=/bin/sh -c "/bin/sleep 15 && /bin/rm -fr /opt/fulcrum/data/ && /usr/local/bin/bitcoin-cli -regtest  -rpcuser=alice -rpcpassword=password -rpcport=18443 generatetoaddress 105 bchreg:qpttdv3qg2usm4nm7talhxhl05mlhms3ys43u76rn0 && /bin/sleep 5 && /opt/fulcrum/Fulcrum-1.2.10-x86_64-linux/Fulcrum /opt/fulcrum/fulcrum.conf"
directory=/opt/fulcrum/
environment=HOME="/opt/fulcrum",PATH="/usr/bin:/usr/sbin"
user=fulcrum
group=fulcrum
autostart=true
autorestart=true
startsecs=60
startretries=3
stopwaitsecs=600
stopsignal=TERM